# systemd/aialra-opencodexapp-runtime-prune.timer
[Unit]
Description=Daily AIALRA Open CodexApp runtime artifact prune

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m
Unit=aialra-opencodexapp-runtime-prune.service

[Install]
WantedBy=timers.target
