# systemd/aialra-opencodexapp-web.service
[Unit]
Description=AIALRA Open CodexApp web bridge
After=network-online.target aialra-opencodexapp-app-server.service
Wants=network-online.target
Requires=aialra-opencodexapp-app-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/srv/aialra/apps/open-codex-app-web-gateway
Environment=HOME=/srv/aialra/state/root-home
Environment=USER=root
Environment=NODE_ENV=production
Environment=CODEX_HOME=/srv/aialra/state/opencodexapp-codex-home
Environment=CODEXAPP_CODEX_CLI=/usr/bin/codex
Environment=CODEXAPP_APP_VERSION=0.137.0
Environment=CODEXAPP_CODEX_PACKAGE_JSON=/usr/lib/node_modules/@openai/codex/package.json
Environment=CODEXAPP_STATE_DIR=/srv/aialra/apps/open-codex-app-web-gateway/state
Environment=CODEXAPP_WEBVIEW_DIR=/opt/codex-desktop/content/webview
Environment=CODEXAPP_EXTERNAL_APP_SERVER=1
Environment=CODEXAPP_WEB_PORT=12920
Environment=CODEXAPP_APP_SERVER_PORT=12921
EnvironmentFile=/srv/aialra/config/secrets/opencodexapp.env
ExecStart=/usr/bin/node /srv/aialra/apps/open-codex-app-web-gateway/web-server.js
Restart=always
RestartSec=3
KillMode=process
StandardOutput=append:/srv/aialra/logs/opencodexapp/web/service.log
StandardError=append:/srv/aialra/logs/opencodexapp/web/service.log
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
